# configs/supervisor.conf
[program:django]
command=gunicorn project.asgi:application -c configs/gunicorn.py -k uvicorn.workers.UvicornWorker
directory=/home/oms
autostart=true
autorestart=true
